# templates/default/drone.conf.erb
start on (filesystem and net-device-up)

chdir /var/lib/drone
console log

env DRONE_TMP=<%= default['drone']['drone_tmp'] %>

script
    DRONED_OPTS="<%= @droned_opts %>"
    if [ -f /etc/default/$UPSTART_JOB ]; then
        . /etc/default/$UPSTART_JOB
    fi
    droned $DRONED_OPTS
end script
